# luci-app-aria2/files/aria2delaystarter.init
#!/bin/sh /etc/rc.common

START=99

start() {
	sleep 10
	/etc/init.d/aria2 start
}
