# luci-app-autoap/luci-app-autoap-ramips-1.0.0/files/etc/init.d/runAutoAP
#!/bin/sh /etc/rc.common

START=99

start() {
        /usr/bin/autoAP &
}

stop() {
        killall -9 autoAP
}

restart() {
        stop
		start
}