# automount/files/automount.init
#!/bin/sh /etc/rc.common
# (C) 2015 openwrt.org

START=45

start() {
	MOUNTED=$(cat '/proc/mounts')'\n'$(cat '/proc/swaps')
	ls /dev/sd*|grep '[a-zA-Z][1-9]'|while read DISK
	do
		MOUNTED_DISK=$(echo -e "$MOUNTED"|grep "^$DISK")
		[ -z "$MOUNTED_DISK" ] && /usr/bin/automount add "${DISK##*/}"
	done
}

stop() {
	/usr/bin/automount umountall
}
